# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:flask-server]
command=python /app/server/main.py
directory=/app/server
autostart=true
autorestart=true
stderr_logfile=/var/log/flask.err.log
stdout_logfile=/var/log/flask.out.log
environment=PYTHONPATH="/app/server"

[program:sveltekit-server]
command=node /app/client/build/index.js
directory=/app/client
autostart=true
autorestart=true
stderr_logfile=/var/log/sveltekit.err.log
stdout_logfile=/var/log/sveltekit.out.log
environment=PORT="3000",HOST="127.0.0.1",PUBLIC_API_URL="http://localhost:5000"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log